# timelapse-off.service
[Unit]
Conflicts=timelapse.service

[Service]
Type=oneshot
ExecStart=/bin/true
