# timelapse.timer
[Timer]
OnCalendar=05..20:0/30

[Install]
WantedBy=timers.target
